# scripts/systemd/docsync-remote-sync.path
[Unit]
Description=DocSync Remote Client — Watch for documentation changes

[Path]
PathModified=%h/docs/
PathModified=%h/README.md
PathModified=%h/CLAUDE.md

[Install]
WantedBy=multi-user.target
